# sync_saves.timer
[Unit]
Description=Backup emudeck saves daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
