# sync_saves.service
[Unit]
Description=Backup emudeck saves

[Service]
ExecStart=/home/deck/scripts/sync_saves.sh
