# sync_saves.timer
[Unit]
Description=Backup emudeck saves hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
